# etc/cron.d/tomsk-realty
# Getting new source data.

PROJECT_DIR=/home/share/www/realty-graph

@daily www-data $PROJECT_DIR/bin/python $PROJECT_DIR/project/src/get_source_data.py new > /var/log/realty-graph.log 2>&1
